# 42.1/Docker/docker-1/02_bonus/lemp/supervisor-lemp.conf
[program:nginx]
command=nginx

[program:php7.0]
command=php7.0

[program:mariadb]
command=mariadb